# supervisord.otp.conf
[program:otp]
command=/etc/freeradius/3.0/otp/otp.rb
autorestart=true
